# supervisord.conf
[supervisord]
nodaemon=true

[program:apache2]
command=/bin/bash -c "source /etc/apache2/envvars && exec /usr/sbin/apache2 -DFOREGROUND"
redirect_stderr=true
autostart=true
autorestart=true

[program:cron]
command = cron -f -L 15
autostart=true
autorestart=true

[program:tail]
command = tail -f /var/log/cron.log
autostart=true
autorestart=true